# ros/kinetic/kinetic-desktop-full/ide/config/supervisord.conf
[supervisord]
nodaemon=false
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:sshd]
command=/usr/sbin/sshd -D
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:chmod]
command=/bin/chmod a+rwx /dev/dri/card0
user=root
autostart=true
stopsignal=QUIT
