# ros/indigo/indigo-desktop-full/amd-intel-graphics/gazebo2/supervisord.conf
[supervisord]
nodaemon=false

[program:sshd]
command=/usr/sbin/sshd -D
user=root
autostart=true
autorestart=true
stopsignal=QUIT

[program:chmod]
command=/bin/chmod a+rwx /dev/dri/card0
user=root
autostart=true
stopsignal=QUIT
